# claude-tracker.service
[Unit]
Description=Claude Usage Optimizer Tracker
After=network.target

[Service]
Type=simple
User=claude-tracker
Group=claude-tracker
WorkingDirectory=/opt/claude-usage-optimizer
EnvironmentFile=/etc/claude-sender.env
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
